# resonaut-server.service
[Unit]
Description=Resonaut WebSocket server
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/resonautweb
ExecStart=npm run server-wss
Restart=on-failure
EnvironmentFile=-/etc/resonaut.env

[Install]
WantedBy=multi-user.target
